# midibox.service
[Unit]
Description=Midibox

[Service]
ExecStart=/home/pi/midibox
User=pi

[Install]
WantedBy=multi-user.target